# deploy/solarx-rollup.timer
[Unit]
Description=Run the SolaX daily_yield rollup nightly

[Timer]
OnCalendar=*-*-* 00:30:00
Persistent=true

[Install]
WantedBy=timers.target
